# .config/systemd/user/git-sync@.path
[Unit] 
Description=Git Sync Path: %I
[Path]
PathModified=%h/%I/.git/refs/heads


[Install]
WantedBy=multi-user.target

