# .config/systemd/user/git-sync@.service
[Unit]
Description=Git Sync: %I
[Service]
Type=oneshot
WorkingDirectory=%h/%I
ExecStart=%h/bin/git-sync
