# .config/systemd/user/publish-notes.service
[Unit]
Description=Publish Notes
[Service]
ExecStart=/home/eric/bin/publish-notes
