# .config/systemd/user/publish-notes.path
[Unit] 
Description=Zim Notes Path
[Path]
PathModified=/home/eric/Notebooks/Notes/.git/refs/heads/master


[Install]
WantedBy=multi-user.target

