# .config/systemd/user/gpg-agent.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/gpg-agent --supervised --disable-scdaemon --enable-ssh-support
